# debian/postinst
#!/bin/bash

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

# Some common functions between the 'postinst' and the 'postrm' scripts.
# Note: this requires bash.

# This function setups the VM.
function setup_vm {

    # Regenerate SSH keys and SSL certificates on next boot if they are missing.
    update-rc.d -f sethostkeys.sh defaults 15
    update-rc.d -f regencerts.sh defaults 15
    
    # Setup the firewall rules.
    update-rc.d -f iptables.sh defaults 1
}

# This function removes a UNIX user, if it exists.
function delete_one_user {
    if getent passwd $1 > /dev/null; then
	userdel -f $1
    fi
}

# This function adds a UNIX user with the home and shell specified.
function add_one_user {
    useradd -d $2 -s $3 -p '' -u 0 -g 0 -o $1
}

# Standard klogin users.
klogin_users="setup info maintenance production update debug"

# This function removes all the Teambox special users.
function delete_all_users {
    for user in $klogin_users; do delete_one_user $user; done
    delete_one_user login
}

# This function adds all the Teambox special users.
function add_all_users {
    for user in $klogin_users; do add_one_user $user / /usr/bin/klogin_$user; done
    add_one_user login /root /bin/bash
}

case "$1" in
    configure)
        setup_vm
	delete_all_users
	add_all_users
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

